# nodejs.conf
[program:nodesprocess]
command=npm run production
directory=/
autostart=true
autorestart=true
startretries=3
